# scripts/com.zion.compost-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zion.compost-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zion/.bun/bin/bun</string>
        <string>packages/compost-cli/src/main.ts</string>
        <string>daemon</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zion/Repos/Zylo/Compost</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/zion/.local/bin:/Users/zion/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/zion</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- KeepAlive=true is safe here (unlike Athena which had IBKR-down crash-loop).
         Compost daemon HC-1 boots even if Engram is unreachable; no external
         hard-dependency that fails closed. ThrottleInterval=60 prevents
         crash-loop if a true bug ships. -->
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Pino writes structured JSON to stdout; launchd appends to daemon.log
         to keep one canonical pino-stream file (matches pre-launchd setup
         where `nohup ... >> daemon.log` was used). -->
    <key>StandardOutPath</key>
    <string>/Users/zion/.compost/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zion/.compost/logs/launchd-stderr.log</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
