# scripts/com.zion.compost-dogfood-7d.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zion.compost-dogfood-7d</string>

    <!-- Runs `claude -p` with the prompt file, appends markdown report
         to ~/.compost/dogfood-7d-report.md. Local TZ (America/Toronto). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>{ echo "# Compost Dogfood 7-Day Report"; echo "Generated: $(date '+%Y-%m-%d %H:%M:%S %Z')"; echo; /Users/zion/.local/bin/claude -p "$(cat /Users/zion/Repos/Zylo/Compost/scripts/dogfood-7d-prompt.txt)"; } &gt;&gt; /Users/zion/.compost/dogfood-7d-report.md 2&gt;&amp;1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zion/Repos/Zylo/Compost</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/zion/.local/bin:/Users/zion/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/zion</string>
    </dict>

    <!-- One-shot fire on 2026-05-04 09:03 local time (America/Toronto).
         Note: launchd's StartCalendarInterval has no Year field, so this
         pattern technically recurs every May 4 at 09:03. The runner prompt
         tells the user to bootout + delete the plist after reading the
         report. If unread for a year, you get a duplicate report -- harmless.
         Minute=3 chosen to avoid :00/:30 fleet alignment. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Month</key>
        <integer>5</integer>
        <key>Day</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>3</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <!-- StandardOutPath/ErrorPath capture any escaped output the script
         missed. Real report goes to ~/.compost/dogfood-7d-report.md via
         the bash redirect inside ProgramArguments. -->
    <key>StandardOutPath</key>
    <string>/Users/zion/.compost/logs/dogfood-7d-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zion/.compost/logs/dogfood-7d-stderr.log</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
